# tests/data/autorestart_test.ini
[unix_http_server]
file=/tmp/test_supervisor.sock

[supervisord]
logfile=/tmp/supervisord_autorestart.log
loglevel=debug
user=root
childlogdir=/tmp

[supervisorctl]
serverurl=unix:///tmp/test_supervisor.sock

[program:crash_test]
command=/bin/false
autorestart=true
startretries=3
startsecs=1
