# tests/data/log_test.ini
[unix_http_server]
file=/tmp/test_supervisor.sock

[supervisord]
logfile=/tmp/supervisord_log_test.log
loglevel=debug
user=root
childlogdir=/tmp

[supervisorctl]
serverurl=unix:///tmp/test_supervisor.sock

[program:counter]
command=/bin/bash -c "for i in {1..10}; do echo Line $i from counter process; sleep 0.1; done"
autorestart=false
stdout_logfile=/tmp/counter.log
stdout_logfile_maxbytes=1KB
redirect_stderr=true

[program:rapid_output]
command=/bin/bash -c "for i in {1..100}; do echo Output line $i; done"
autorestart=false
stdout_logfile=/tmp/rapid_output.log
stdout_logfile_maxbytes=500
redirect_stderr=true
