# tests/data/test_config.ini
[unix_http_server]
file=/tmp/test_supervisor.sock

[supervisord]
logfile=/tmp/test_supervisord.log
loglevel=debug
user=testuser
childlogdir=/tmp/test_logs

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/test_supervisor.sock

[program:test_app]
command=/bin/echo test_app
environment=FOO=bar,BAZ=qux
directory=/tmp
autorestart=true
user=root
stdout_logfile=/tmp/%(program_name)s.log
stdout_logfile_maxbytes=5MB
redirect_stderr=true
startsecs=2
startretries=5
stopwaitsecs=15
stopsignal=TERM
